# restarter@.service
[Unit]
Description=perform restart a service

[Service]
ExecStart=/bin/systemctl --user restart %i
Type=oneshot

[Install]
WantedBy=default.target
